# Library/LaunchAgents/me.kosui.ollama-keep-alive.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>me.kosui.ollama-keep-alive</string>
  <!-- ログインごとに OLLAMA_KEEP_ALIVE を GUI セッションへ設定し、
       Ollama.app が起動時に継承してモデルを常駐させる。
       launchctl setenv は再起動で消えるため LaunchAgent で永続化する。 -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/launchctl</string>
    <string>setenv</string>
    <string>OLLAMA_KEEP_ALIVE</string>
    <string>-1</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
